# d5next/d5next.service
[Unit]
Description=d5next job

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/d5next.sh
